# system/systemd/apt-daily.timer.d/vibes.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 22:00
RandomizedDelaySec=15m
